# opscenter/src/supervisord.conf
[supervisord]
nodaemon=true

[program:opscenter]
user=root
command=/usr/share/opscenter/bin/opscenter
